# systemd/user/pipewire.service.d/override.conf
[Service]
MemoryDenyWriteExecute=no
ExecStartPost=setsid -f env -C /home/missed/.local/share/java8 java Main
